# systemd/vllm-qwen3.6-35b-a3b-dflash.service
[Unit]
Description=vLLM Qwen3.6-35B-A3B-DFlash
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=%h/dgx-spark-vllm-qwen3.6-35b-a3b-dflash
ExecStart=%h/dgx-spark-vllm-qwen3.6-35b-a3b-dflash/scripts/start.sh
ExecStop=-docker stop vllm-qwen3.6-35b-a3b-dflash
ExecStop=-docker rm vllm-qwen3.6-35b-a3b-dflash

[Install]
WantedBy=default.target
